# systemd/deadlock-tracker.service.example
[Unit]
Description=Deadlock SteamDB Updates Tracker Bot
After=network.target

[Service]
Type=simple
# Вкажіть користувача, від імені якого буде запускатися бот
User=root

# Шлях до робочої директорії бота
WorkingDirectory=/opt/deadlock-tracker

# Шлях до Python у віртуальному оточенні (або глобального python3) та файлу bot.py
ExecStart=/opt/deadlock-tracker/venv/bin/python bot.py

# Автоматичний перезапуск у разі падіння
Restart=always
RestartSec=10

# Логування стандартного виводу
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=deadlock-tracker

[Install]
WantedBy=multi-user.target
